# docker/php/crontab
* * * * * php /app/src/scripts/update_device_statuses.php >> /var/log/cron.log 2>&1